# frontend/timer_frontend.service
[Unit]
Description=Slot Car Timer Frontend
After=slot_timer.service timer_update.service network-online.target disable_cursor.service
Wants=slot_timer.service timer_update.service network-online.target disable_cursor.service

[Service]
WorkingDirectory=/home/pi/SlotRacr/frontend
ExecStart=/home/pi/SlotRacr/frontend/app.py
User=pi
StandardInput=tty
StandardOutput=tty
TTYPath=/dev/tty2
TimeoutSec=10

[Install]
WantedBy=multi-user.target