# BackEnd/slot_timer.service
[Unit]
Description=Slot Car Timer Backend

[Service]
ExecStart=/usr/bin/slotcar.py

[Install]
WantedBy=multi-user.target